# etc/systemd/system/streamer-webapp.service.d/safeguards.conf
[Service]
MemoryMax=150M
MemoryHigh=120M
Restart=always
RestartSec=5
